# systemd/user/calendarsync.service
[Unit]
Description=calendarsync runner

[Service]
Type=oneshot
ExecStart=/home/jmickey/bin/calendarsync
TimeoutSec=600s
