# systemd/user/calendarsync.timer
[Unit]
Description=Sync my calendar once every 10 minutes

[Timer]
OnBootSec=60
OnUnitActiveSec=600

[Install]
WantedBy=timers.target