# systemd/user/main.target
[Unit]
Description=Main Window Manager Target
Requires=default.target
