# systemd/user/offlineimap.service
[Unit]
Description=OfflineImap Runner

[Service]
ExecStart=/home/codemac/bin/mailupd
TimeoutSec=50
