# systemd/user/offlineimap.service
[Unit]
Description=OfflineImap Runner

[Service]
Type=oneshot
ExecStart=/home/codemac/bin/mailupd
TimeoutSec=600s
