# systemd/user/offlineimap.service
[Unit]
Description=OfflineImap Runner

[Service]
Type=oneshot
ExecStart=/home/jmickey/bin/mailupd
TimeoutSec=600s
