# systemd/user/offlineimap.timer
[Unit]
Description=OfflineImap Every 1m

[Timer]
OnBootSec=60
OnUnitActiveSec=60

[Install]
WantedBy=timers.target