# systemd/user/pulseaudio.service
[Unit]
Description=PulseAudio Sound System

[Service]
ExecStart=/usr/bin/pulseaudio

[Install]
WantedBy=default.target
