# systemd/user/shepherd.service
[Unit]
Description=Shepherd herds the services

[Service]
ExecStart=/bin/shepherd
EnvironmentFile=/home/jmickey/config/env-arch