# systemd/user/xbanish.service
[Unit]
Description=xbanish

[Service]
ExecStart=/usr/bin/xbanish
Restart=always
EnvironmentFile=/home/codemac/config/env-guix

[Install]
WantedBy=main.target
